# ansible/playbooks/roles/auto_deploy/templates/ansible-deploy.timer.j2
[Unit]
Description=Ansible Playbook Auto-Deploy Timer
Requires=ansible-deploy.service

[Timer]
OnCalendar=Tue *-*-* 04:00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
