# ansible/playbooks/roles/wealth_tracker/templates/wealth_tracker_update.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Run Wealth Tracker update periodically

[Timer]
OnCalendar=03:00

[Install]
WantedBy=timers.target
