# ansible/playbooks/roles/wealth_tracker/templates/wealth_tracker_update.service.j2
# {{ ansible_managed }}

[Unit]
Description=Update Wealth Tracker container

[Service]
Type=oneshot
ExecStart=/usr/local/bin/deploy_wealth_tracker.sh
