# ansible/playbooks/roles/docker/templates/docker-prune.timer.j2
[Unit]
Description=Weekly Docker Image Prune Timer
Requires=docker-prune.service

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
